# .crontab
# m h dom mon dow command

@reboot $HOME/cron/init.sh

# stat
00 09 * * * $HOME/cron/stat.sh "$HOME/Downloads/1091/stat" log

# ぐーもる荘
* * * * * for i in $(seq 0 20 59); do (sleep $i; $HOME/cron/gmr-config-sync.sh log) & done
* * * * * for i in $(seq 3 20 59); do (sleep $i; $HOME/cron/gmr-check-feeds.sh "$HOME/Downloads/1091/gu-morusou/feeds-owner" fc "$HOME/Downloads/1091/gu-morusou/.lock_check-feeds-owner") & done
* * * * * for i in $(seq 5 20 59); do (sleep $i; $HOME/cron/gmr-check-likes.sh "$HOME/Downloads/1091/gu-morusou/like_latest" fc "$HOME/Downloads/1091/gu-morusou/.lock_check-likes") & done
* * * * * for i in $(seq 8 20 59); do (sleep $i; $HOME/cron/gmr-feed-fetch.sh "$HOME/Downloads/1091/gu-morusou/feeds/" 10 "$HOME/Downloads/1091/gu-morusou/.lock_feeds-fetch") & done
* * * * * for i in $(seq 8 20 59); do (sleep $i; $HOME/cron/gmr-activity-fetch.sh "$HOME/Downloads/1091/gu-morusou/activity/" 20 "$HOME/Downloads/1091/gu-morusou/.lock_activity-fetch") & done
00 05 * * * $HOME/cron/gmr-hitokoto-download.sh "$HOME/Downloads/1091/gu-morusou/hitokoto" log
02 05 * * * $HOME/cron/gmr-hitokoto-upload.sh $HOME/Downloads/1091/gu-morusou/hitokoto/ gdrive:1091/resources/fc/gu-morusou.app/hitokoto/ log
05 05 28 * * $HOME/cron/gmr-hitokoto-clean.sh $HOME/Downloads/1091/gu-morusou/hitokoto/ gdrive:1091/resources/fc/gu-morusou.app/hitokoto/ log
* * * * * $HOME/cron/gmr-check-error.sh "$HOME/Downloads/1091/gu-morusou/error" log

00 00 * * * rm -rf /tmp/tmp.gmr*

# Check AGQR
* * * * * $HOME/cron/record-agqr-check.sh "$HOME/Downloads/agrec-check.err" log

# 徳井青空のまぁるくなぁれ！ (毎週木曜25:00-25:30)
#00 01 * * 5 $HOME/cron/record-agqr-marunare.sh "$HOME/Downloads/1091/marunare" log
# 徳井青空のまぁるくなぁれ！ アーカイブ (毎週金曜13:00)
#05 12 * * 5 $HOME/cron/download-marunare.sh "$HOME/Downloads/1091/marunare" log

# アルマギア情報局 
30 16 * * 0 $HOME/cron/record-agqr-armagia.sh "$HOME/Downloads/1091/agqr" log

# 週刊ファルコムラジオ (毎週金曜22:00-22:14)
00 22 * * 5 $HOME/cron/record-agqr-falcom.sh "$HOME/Downloads/1091/agqr" log

